# resources/swap-me.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/swap-me

[Install]
WantedBy=multi-user.target
